# dcache_operations/files/etc/profile.d/java.sh
#!/bin/bash

# Managed by Puppet (dcache)
# Do NOT edit, changes will be overwritten!

export JAVA_HOME=/usr/java/default
export PATH=$PATH:/usr/java/default/bin
